# mjpg-streamer/mjpg-streamer-experimental/mjpg-streamer.service
[Unit]
Description=MJPG Streamer service

[Service]
ExecStart=${MJPG_STREAMER_FILE_PATH}